# scripts/claude-bot.service
[Unit]
Description=Claude Bot — tmux sessions + admin daemon
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=forking
User=%i
WorkingDirectory=/home/%i/bots/claude-bot
ExecStartPre=/home/%i/.local/bin/claude-bot up
ExecStart=/usr/bin/env bun /home/%i/bots/claude-bot/scripts/admin-daemon.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
